# deploy/systemd/omnitun-server.service
[Unit]
Description=OmniTun API Server
Documentation=https://docs.omnitun.io
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=omnitun
Group=omnitun
ExecStart=/usr/local/bin/omnitun-server --config /etc/omnitun/config.yaml
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/omnitun /var/log/omnitun
PrivateTmp=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Observability
Environment=OTEL_SERVICE_NAME=omnitun-server

[Install]
WantedBy=multi-user.target
